# etc/init/flower-server-service.conf
# Should be /etc/init/redis-server.conf
# If redis-server is started via init.d, it must be removed from there.
#    sudo update-rc.d -f redis-server remove

description "Start the Flower: Celery Monitoring Panel"
author "Yash <yash@abc.com>"

start on runlevel[2345]

stop on runlevel[06]

respawn

respawn limit 3 30

chdir /home/ec2-user/www/ABCBackend/abc

script
    source /home/ec2-user/www/ABCBackendEnv/bin/activate
    # celery flower -A abc --address=127.0.0.1 --port=5555
    flower -A abc  --port=5555 --broker=redis://localhost:6379/0
    # flower --app=abc --port=5555 --broker=redis://localhost:6379/0
end script


respawn
